# .config/systemd/user/fix-edge-wmclass.service
[Unit]
Description=Fix StartupWMClass in Microsoft Edge .desktop files

[Service]
Type=oneshot
ExecStart=%h/.local/bin/fix-edge-wmclass
