# .config/systemd/user/fix-edge-wmclass.path
[Unit]
Description=Watch applications dir for Edge .desktop changes

[Path]
PathModified=%h/.local/share/applications

[Install]
WantedBy=default.target
