# playbooks/system/roles/base/files/rc.local
#!/bin/sh -e

if [ -x "$(command -v teamviewer)" ]; then
    teamviewer --daemon disable
fi
